# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/app/storage/logs/supervisord.log
pidfile=/app/storage/supervisord.pid
loglevel=info

[program:php-fpm]
command=php-fpm -F
directory=/app
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/app/storage/logs/php-fpm.log

[program:caddy]
command=caddy run --config /app/Caddyfile --adapter caddyfile
directory=/app
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/app/storage/logs/caddy.log

[program:laravel-worker]
command=php artisan queue:work --tries=3 --timeout=300 --verbose
directory=/app
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/app/storage/logs/worker.log
process_name=%(program_name)s_%(process_num)02d
numprocs=1